# agent_intent_test.service
[Unit]
Description=Agent Intent Test Framework
After=network.target

[Service]
Type=simple

User=agent
Group=agent

WorkingDirectory=/opt/agent_intent_test_framework

ExecStart=/usr/bin/python3 \
  /opt/agent_intent_test_framework/manage.py \
  runserver 0.0.0.0:18081

Restart=always
RestartSec=5

StandardOutput=append:/opt/agent_intent_test_framework/logs/server.log
StandardError=append:/opt/agent_intent_test_framework/logs/server.log

[Install]
WantedBy=multi-user.target
